# dist/supervisord.conf
[supervisord]
nodaemon=true
identifier=supervisord

[include]
files=run/supervisord.conf

[program:libsrv]
command=libctl.sh libsrv
user=nobody
group=nogroup

[program:httpd]
command=libctl.sh httpd
