# dist/supervisord.conf
[supervisord]
nodaemon=true
identifier=supervisord

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[include]
files=var/*.supervisor

[program:libsrv]
command=libctl.sh libsrv
[program:httpd]
command=libctl.sh httpd
