# dist/supervisord.conf
[supervisord]
nodaemon=true
identifier=supervisord

[include]
files=var/*.supervisor

[program:libsrv]
command=libctl.sh libsrv
[program:httpd]
command=libctl.sh httpd
